# config/includes.chroot/etc/xdg/autostart/nm-applet.desktop
[Desktop Entry]
Name=NetworkManager Applet
Exec=nm-applet
Type=Application
X-GNOME-Autostart-enabled=true
